# wallpapers/.config/systemd/user/change-wallpaper.timer
[Unit]
Description=Change wallpaper every 30 mins

[Timer]
Persistent=true
OnUnitActiveSec=30min
OnBootSec=10min

[Install]
WantedBy=timers.target
